# scripts/chzzk-recorder.service
[Unit]
Description=Chzzk-Recorder-Pro - Chzzk Live Recording Service
Documentation=https://github.com/your-username/Chzzk_downloader
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# 실행 사용자 (root 외 별도 사용자 권장)
User=%i
Group=%i

# 작업 디렉토리: 프로젝트 루트로 변경 필요
WorkingDirectory=/opt/chzzk-recorder/backend
ExecStart=/opt/chzzk-recorder/.venv/bin/python run.py

# 환경 변수
EnvironmentFile=/opt/chzzk-recorder/.env

# 재시작 정책
Restart=on-failure
RestartSec=5s
StartLimitInterval=60s
StartLimitBurst=3

# 로그 설정
StandardOutput=journal
StandardError=journal
SyslogIdentifier=chzzk-recorder

# 보안 강화 (선택사항)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
